# deploy/systemd/7shifts-birthday-calendar.service
[Unit]
Description=7shifts Birthday Calendar API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/7shifts-birthday-calendar
EnvironmentFile=/etc/7shifts-birthday-calendar.env
ExecStart=/usr/bin/node /opt/7shifts-birthday-calendar/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
